# src/digi_office/deploy/digi-office.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>digi-office.agent</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>$HOME/.config/digi-office/scripts/ciphemon_agent.py</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>$HOME/.config/digi-office/logs/digi-office.log</string>
  <key>StandardErrorPath</key>
  <string>$HOME/.config/digi-office/logs/digi-office.err</string>
  <key>WorkingDirectory</key>
  <string>$HOME/.config/digi-office/workspace/digi_office</string>
</dict>
</plist>
